# src/django_mosaic/conf/deployment/mosaic-backup.timer
[Unit]
Description=Hourly backup of {{APP_NAME}} database
Requires=mosaic-backup.service

[Timer]
# Run every hour on the hour
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
